# scheduler.cron
* * * * * php /home/container/webroot/artisan schedule:run >> /home/container/logs/cron.log 2>&1
