# scheduler.cron
* * * * * php /home/container/webroot/artisan schedule:run >> /dev/null 2>&1
